# pi/scripts/deltaquest_encryption/linux_server_end/stunnel-deltaquest-[siteid]-[portid].service
[Unit]
Description="DeltaQuest Stunnel Tunnel - siteid portid"
After=network.target

[Service]
ExecStart=stunnel /u/tunnels/deltaquest-siteid-portid-tunnel
RestartSec=15
Restart=always
KillMode=mixed

[Install]
WantedBy=multi-user.target
