# socat/files/service_files/rhel8/socat10.service
[Unit]
Description='Socat Access Control Service 10'
After=network.target

[Service]
Type=simple
User=root
ExecStart=/bin/bash -C '/u/espbin/socat/startsocat_scripts/startsocat10'
RestartSec=10
Restart=on-failure
KillMode=process
PIDFile=/var/run/socat10.pid

[Install]
WantedBy=multi-user.target
