# pi/scripts/http_interface_secure/linux_server_end/stunnel-httpsecure-[siteid]-[eliteunithttpport].service
[Unit]
Description="http Interface Stunnel Tunnel - siteid eliteunithttpport"
After=network.target

[Service]
ExecStart=/usr/bin/stunnel /u/tunnels/httpsecure-siteid-eliteunithttpport-tunnel
RestartSec=15
Restart=always
KillMode=mixed

[Install]
WantedBy=multi-user.target
